# ez-volume.service
[Unit]
Description=EZ Volume Service

[Service]
ExecStart=/usr/local/bin/ez-volume.sh

[Install]
WantedBy=graphical-session.target

